# storage/rootfs/etc/s6-overlay/s6-rc.d/storage/run
#!/command/with-contenv bashio

bashio::log.info "Starting Storage nginx on port 8099..."

MY_HOSTNAME="${HOSTNAME:-}"
[ -z "${MY_HOSTNAME}" ] && MY_HOSTNAME=$(hostname 2>/dev/null) || true
MY_HOSTNAME_NORM=$(echo "${MY_HOSTNAME:-}" | tr '_' '-')
REPO_PREFIX="${MY_HOSTNAME_NORM%%-ha-storage}"

_detect_scraper() {
    local _url=""
    if [ -n "${REPO_PREFIX}" ] && [ "${REPO_PREFIX}" != "${MY_HOSTNAME_NORM}" ]; then
        for _c in "${REPO_PREFIX}-grocy-scraper" "${REPO_PREFIX}_grocy_scraper"; do
            if curl -sf --max-time 3 "http://${_c}:8099/api/config" >/dev/null 2>&1; then
                _url="http://${_c}:8099"
                break
            fi
        done
    fi
    if [ -z "${_url}" ] && [ -n "${SUPERVISOR_TOKEN}" ]; then
        local _slug=""
        _slug=$(curl -sf --max-time 5 \
            -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
            http://supervisor/addons \
            | jq -r '[.data.addons[]? | select(.slug | test("grocy.scraper"; "i"))][0].slug // empty') || true
        if [ -n "${_slug}" ]; then
            local _ip=""
            _ip=$(curl -sf --max-time 5 \
                -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
                "http://supervisor/addons/${_slug}/info" \
                | jq -r '.data.ip_address // empty') || true
            if [ -n "${_ip}" ]; then
                _url="http://${_ip}:8099"
            elif curl -sf --max-time 3 "http://${_slug}:8099/api/config" >/dev/null 2>&1; then
                _url="http://${_slug}:8099"
            fi
        fi
    fi
    echo "${_url}"
}

_SCRAPER_PLACEHOLDER="http://127.0.0.1:1"
SCRAPER_AUTODETECT=0
SCRAPER_URL=""

if bashio::config.has_value 'scraper_url'; then
    SCRAPER_URL=$(bashio::config 'scraper_url')
    bashio::log.info "Scraper URL (from config): ${SCRAPER_URL}"
else
    SCRAPER_AUTODETECT=1
    bashio::log.info "Attempting auto-detect of Grocy Scraper addon (hostname: '${MY_HOSTNAME}')..."
    SCRAPER_URL=$(_detect_scraper)
    if [ -n "${SCRAPER_URL}" ]; then
        bashio::log.info "Scraper URL (auto-detected): ${SCRAPER_URL}"
    else
        bashio::log.info "Grocy Scraper addon not found -- will retry in background."
    fi
fi

[ -z "${SCRAPER_URL}" ] && SCRAPER_URL="${_SCRAPER_PLACEHOLDER}"
export SCRAPER_URL

envsubst '${SCRAPER_URL}' < /etc/nginx/nginx.conf.template > /etc/nginx/nginx.conf

for i in $(seq 1 30); do
    if curl -sf http://127.0.0.1:8100/api/health > /dev/null 2>&1; then
        bashio::log.info "Storage API is ready."
        break
    fi
    bashio::log.info "Waiting for Storage API (attempt ${i}/30)..."
    sleep 2
done

bashio::log.info "Scraper URL: ${SCRAPER_URL}"

# Background health monitor: probe every 30 s when in auto-detect mode.
# - If not found: keep trying.
# - If found but unreachable: re-detect; reload nginx only if URL changed or reverted.
if [ "${SCRAPER_AUTODETECT}" = "1" ]; then
    bashio::log.info "Starting background scraper health monitor (every 30 s)..."
    (
        _cur="${SCRAPER_URL}"
        while sleep 30; do
            if [ "${_cur}" = "${_SCRAPER_PLACEHOLDER}" ]; then
                _new=$(_detect_scraper)
                if [ -n "${_new}" ]; then
                    bashio::log.info "Scraper found at ${_new} -- reloading nginx."
                    _cur="${_new}"
                    SCRAPER_URL="${_cur}"
                    envsubst '${SCRAPER_URL}' < /etc/nginx/nginx.conf.template > /etc/nginx/nginx.conf
                    nginx -s reload
                fi
            elif ! curl -sf --max-time 3 "${_cur}/api/config" >/dev/null 2>&1; then
                bashio::log.warning "Scraper at ${_cur} unreachable -- re-detecting..."
                _new=$(_detect_scraper)
                if [ -n "${_new}" ] && [ "${_new}" != "${_cur}" ]; then
                    bashio::log.info "Scraper re-detected at ${_new} -- reloading nginx."
                    _cur="${_new}"
                    SCRAPER_URL="${_cur}"
                    envsubst '${SCRAPER_URL}' < /etc/nginx/nginx.conf.template > /etc/nginx/nginx.conf
                    nginx -s reload
                elif [ -z "${_new}" ]; then
                    bashio::log.warning "Scraper not found -- reverting to placeholder."
                    _cur="${_SCRAPER_PLACEHOLDER}"
                    SCRAPER_URL="${_cur}"
                    envsubst '${SCRAPER_URL}' < /etc/nginx/nginx.conf.template > /etc/nginx/nginx.conf
                    nginx -s reload
                fi
            fi
        done
    ) &
fi

exec nginx -g "daemon off;"
