# storage/rootfs/etc/s6-overlay/s6-rc.d/storage/run
#!/command/with-contenv bashio

bashio::log.info "Starting Storage nginx on port 8099..."

# Copy nginx template (no envsubst needed — static config)
cp /etc/nginx/nginx.conf.template /etc/nginx/nginx.conf

# Wait for FastAPI to be ready
for i in $(seq 1 30); do
    if curl -sf http://127.0.0.1:8100/api/health > /dev/null 2>&1; then
        bashio::log.info "Storage API is ready."
        break
    fi
    bashio::log.info "Waiting for Storage API (attempt ${i}/30)..."
    sleep 2
done

exec nginx -g "daemon off;"
