# storage/rootfs/etc/s6-overlay/s6-rc.d/storage/run
#!/command/with-contenv bashio

bashio::log.info "Starting HA-Storage addon..."

# Read config
GEMINI_API_KEY=$(bashio::config 'gemini_api_key')
GEMINI_MODEL=$(bashio::config 'gemini_model')
DEBUG=$(bashio::config 'debug')

export GEMINI_API_KEY
export GEMINI_MODEL
export DEBUG

# Data directory for SQLite DB and images (persisted across restarts)
export DATA_DIR="/data"
mkdir -p "${DATA_DIR}/images/products" "${DATA_DIR}/images/recipes"

bashio::log.info "Database: ${DATA_DIR}/storage.db"

exec python3 /app/main.py
